# var/spool/cron/crontabs/root
# min  hour  day  month  weekday  command

10 12 * * * /sbin/acme-client-update myhost.sample.test
